# deploy/trading-galaxy.service
[Unit]
Description=Trading Galaxy KB API
Documentation=https://github.com/Alby2007/Mem-Backend
After=network-online.target ollama.service
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User=__USER__
Group=__USER__
WorkingDirectory=__APP_DIR__
EnvironmentFile=__APP_DIR__/.env
ExecStart=__VENV_DIR__/bin/gunicorn --config __APP_DIR__/gunicorn.conf.py api_v2:app
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=trading-galaxy

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=__DATA_DIR__ __APP_DIR__

[Install]
WantedBy=multi-user.target
